# deploy/maxbot-dashboard.service
[Unit]
Description=MAX BOT - Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/maxbot
ExecStart=/opt/maxbot/venv/bin/gunicorn app:app --bind 0.0.0.0:5001 --timeout 120
Restart=always
RestartSec=10
StandardOutput=append:/opt/maxbot/logs/dashboard.log
StandardError=append:/opt/maxbot/logs/dashboard.log

[Install]
WantedBy=multi-user.target
